# systemd/.config/systemd/user/package-updates.service
[Unit]
Description=Check Package Updates

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/scripts/package-updates.sh recheck
